# .config/dotbare-cfg/wsl-mount-shared.service
[Unit]
Description=Mount Windows network share
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wsl-mount-shared.sh
Restart=no
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
